# deploy/systemd/llmport-agent.service
[Unit]
Description=llm-port node agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=llm-port-agent
Group=llm-port-agent
SupplementaryGroups=docker
EnvironmentFile=-/etc/llmport-agent.env
ExecStart=/usr/local/bin/llmport-agent
Restart=always
RestartSec=5
TimeoutStopSec=30

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=true
ReadWritePaths=/var/lib/llmport-agent
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
